# bloop/native-bin/bloop.socket
[Unit] 
Description = Bloop Scala build server

[Socket] 
ListenStream = 8212

[Install] 
WantedBy = sockets.target